# project.conf
description "uWSGI server instance configured to serve PROJECT"

start on runlevel [2345]

stop on runlevel [!2345]

setuid root

setgid www-data

env PATH=/home/ubuntu/venv/bin

chdir /home/ubuntu/PROJECT

exec uwsgi --ini PROJECT.ini
